# supervisord.conf.d/apache2.conf
[program:apache2]
command=/usr/sbin/apache2ctl -D FOREGROUND
environment=APACHE_LOG_DIR="/opt/aptrepo/log"
stopasgroup=true
