# supervisord.conf.d/rsyslogd.conf
[program:rsyslogd]
command=/usr/sbin/rsyslogd -n
